# examples/scheduled/launchd.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.signalharness.scan</string>
  <key>WorkingDirectory</key>
  <string>/path/to/SignalHarness</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/uv</string>
    <string>run</string>
    <string>signal-harness</string>
    <string>scan</string>
    <string>--fixture</string>
    <string>examples/signal_harness/sample_events.json</string>
    <string>--mode</string>
    <string>mock-agent</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/signalharness.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/signalharness.err.log</string>
</dict>
</plist>

